# cron.d/wordpress_data_weekly_backup
# /etc/cron.d/wordpress_data_weekly_backup: crontab fragment for creating wordpress data backup
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
PROJECT_TAG=test
TARGET_PATH=/backups/data/weekly

# creating wordpress data backup
# starts every week at 02:00 AM on Saturday
0 2 * * 6 root export DATE=$(date +\%H-\%M_\%m-\%d-\%y) && mkdir -p $TARGET_PATH && docker cp $PROJECT_TAG-wordpress:/var/www/html/wp-content $TARGET_PATH/backup_$DATE && zip -qr $TARGET_PATH/backup_$DATE.zip $TARGET_PATH/backup_$DATE && rm -rfd $TARGET_PATH/backup_$DATE

# delete files older than 60 days
# starts every week at 03:00 AM on Saturday
0 3 * * 6 root find $TARGET_PATH -name "backup_*.zip" -type f -mtime +60 -delete
